# service/maild_web.service
[Unit]
Description=Maild Web
#After=network.target   
Requires=postgresql.service
After=postgresql.service


[Service]
Type=simple
ExecStart=/usr/bin/maild_web -c /etc/maild/maild_web.conf -l /etc/maild/web_log.properties
User=maild
Group=maild
Restart=on-abort

[Install]
WantedBy=multi-user.target default.target

